# kms/derive-server-enclave/supervisord.conf
[supervisord]
loglevel=debug
logfile=/dev/stdout
logfile_maxbytes=0

# attestation server
[program:attestation-server]
command=/app/attestation-server --ip-addr 0.0.0.0:1300 --pub-key /app/secp256k1.pub --user-data /app/init-params
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# attestation server x25519
[program:attestation-server-x25519]
command=/app/attestation-server --ip-addr 0.0.0.0:1301 --pub-key /app/x25519.pub --user-data /app/init-params
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# proxy components inside enclave
[program:ip-to-vsock-raw-outgoing]
command=/app/ip-to-vsock-raw-outgoing --vsock-addr 3:1200 --queue-num 0
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

[program:vsock-to-ip-raw-incoming]
command=/app/vsock-to-ip-raw-incoming --vsock-addr 88:1200 --device lo
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# DNS-over-HTTPS provider
[program:dnsproxy]
command=/app/dnsproxy -u https://1.1.1.1/dns-query
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# your custom programs go here

# kms derive server
# WARN: This is only a test server
# You would never want this to listen on 0.0.0.0 in prod
[program:kms-derive-server]
command=/app/kms-derive-server --kms-endpoint v1.kms.box:1100 --listen-addr 0.0.0.0:1100 --pcr0 cb7ebc13d527e9cf9cc271b0d816c72a5bfa685ae56118ce4986fa82b8b9aac8b851206539a6c4600ad77566aa15bd0e --pcr1 3b26340a10ac3a5494139fed12bc30028017b72dabfb78a38763bd21ea67bbfe03214e7ce628c2952cfa2ff478f370ba --pcr2 0d128997bfc8ef24a2aa1ecec60c61e48eda9b439f8662d41eb38d8f0ba0401367778a5f8804127c591d824d5c3ed617 --user-data 544d4b69000100009293c4308e0ebc7f830be9a963f343358c4875d67c36eb1dc0c3250cb9438a43f0399ed4fc7e0b80ccb99b106136e06bbaecb90bc460843d33b632beb3537a7d1c42531ad1b43bce0015ba7e7fb7ed82fed9a0d200ac11dc646ca22b743daec9baa39d689c350a3e70ce85979e856e9cac3c31f45febe28510360beb78495932a6e06c2a5fd6caa0adc9dff0fef20d9ff4916199bcbcc450cf36031706fbe58b82fcf808e76f801924ec66a3ea5bcaec0fb33619cf5df6314b43eca506a25392b1e042126a63744f35bd39403c6d73e81a4efaf0e0251563ee3eacc2ce419ae5a1c139d5d321761a9292c430868c3d012a5d524f0939e4ee4d60b738b4c44448ec286a5361e15ffbf2641e2df25363a204a738231e5f1a9621999741da01b87b22636f6e646974696f6e223a7b22636861696e223a312c22636f6e646974696f6e54797065223a22636f6e7472616374222c22636f6e747261637441646472657373223a22307843374430383443326536424341633030374433424146433431353438356133383038303343306265222c2266756e6374696f6e416269223a7b22696e70757473223a5b7b22696e7465726e616c54797065223a2261646472657373222c226e616d65223a22222c2274797065223a2261646472657373227d5d2c226e616d65223a2269735665726966696564222c226f757470757473223a5b7b22696e7465726e616c54797065223a22626f6f6c222c226e616d65223a22222c2274797065223a22626f6f6c227d5d2c2273746174654d75746162696c697479223a2276696577222c2274797065223a2266756e6374696f6e227d2c226d6574686f64223a2269735665726966696564222c22706172616d6574657273223a5b223a7573657241646472657373225d2c2272657475726e56616c756554657374223a7b22636f6d70617261746f72223a223d3d222c2276616c7565223a747275657d7d2c2276657273696f6e223a22312e302e30227dc441cd0ad8db7d03074a07437a39025a8deea6a72c433dedc45375b95d28257dc0300ae70507050f75c5d5c31c50b167bdde4c06ac789b22641defe91de9442fd74d1b
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0
