# kms/derive-server-enclave/supervisord.conf
[supervisord]
loglevel=debug
logfile=/dev/stdout
logfile_maxbytes=0

# attestation server
[program:attestation-server]
command=/app/attestation-server --ip-addr 0.0.0.0:1300 --pub-key /app/secp256k1.pub --init-params /app/init-params
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# attestation server x25519
[program:attestation-server-x25519]
command=/app/attestation-server --ip-addr 0.0.0.0:1301 --pub-key /app/x25519.pub --init-params /app/init-params
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# proxy components inside enclave
[program:ip-to-vsock-raw-outgoing]
command=/app/ip-to-vsock-raw-outgoing --vsock-addr 3:1200 --queue-num 0
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

[program:vsock-to-ip-raw-incoming]
command=/app/vsock-to-ip-raw-incoming --vsock-addr 88:1200 --device lo
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# DNS-over-HTTPS provider
[program:dnsproxy]
command=/app/dnsproxy -u https://1.1.1.1/dns-query
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# your custom programs go here

# kms derive server
# WARN: This is only a test server
# You would never want this to listen on 0.0.0.0 in prod
[program:kms-derive-server]
command=/app/kms-derive-server --kms-endpoint 35.154.162.233:1100" --listen-addr 0.0.0.0:1100 --pcr0 1a96e41349731a2d234cbf1558f8b8a27d6d5afe169ee1672ea4582b4e4ce5d113116308eb10a0a283515b0989cd4c16 --pcr1 8384311b93170c6ae8b29b927eb683022a118957d0ad26d2c866fdb87e96bc1fcad5583ead068a103189f0c1731afb84 --pcr2 ec7bec94f08fcb51f6c086a8d3cac309902a7f56b573dc44d333b05cdc707948e75bd03ea17bf645245b047d40626071 --user-data 544d4b69000100009293c43099822c1c52b001c0af30a3149886fb9e44f0877b87b0afdb132e34f9fc8d47af91c6dec850e58df5e527d00817ed150cc4608f40d01bd170ce39056de80e9b36d0d70002054d7854119c04389f9030cbc798973e8deab07300108828350727b1ca180d78e21f7b9ca6caf59300c89897a84c34e0b4e623b83e82eba6995b69cca10d2e7b06d4988b54a6d9b26813ae83624fc4503e76999a10241c261a1ec7c18ad92d1c97eae82459090e63fee50b7b51f4898444ec096ec44a4cd0e060081464109766e3bad397a425a2f648d75cd55da33cce3fa08f69946dbe3d89443d554b6874f89292c4308a685c7f7348574a4889ade49e0c57e67995d1d5f588478d2e6457eaf193680731786c37c48192d9931793305b2b5610d9b67b22636f6e646974696f6e223a7b22636861696e223a38303030322c22636f6e646974696f6e54797065223a22727063222c226d6574686f64223a226574685f67657442616c616e6365222c22706172616d6574657273223a5b223a7573657241646472657373222c226c6174657374225d2c2272657475726e56616c756554657374223a7b22636f6d70617261746f72223a223e3d222c2276616c7565223a307d7d2c2276657273696f6e223a22312e302e30227dc44143cf33150a8482e77ddc3dbb8c0ce71d0101c072bc51b7cee74694c8300c4dc3582d1c6159cb066fa00c8777df68937505d0d813ec3317e804c3d4de6b74e2421b36
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0
